# supervisord/3.3.0/supervisord.conf
[unix_http_server]
file=/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://supervisor/supervisor.sock

[supervisord]
logfile=/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/supervisor/supervisord.pid
childlogdir=/supervisor
nodaemon=true

[program:program_0]
command=/bin/sh /supervisor/program_0.sh
autostart=true
autorestart=false

[program:program_1]
command=/bin/sh /supervisor/program_1.sh
autostart=true
autorestart=false

[program:program_2]
command=/bin/sh /supervisor/program_2.sh
autostart=true
autorestart=false
